# systemd/yolobot.service
[Unit]
Description=YOLObot
# systemd user services will start on any host by default, making duplicate bots
ConditionHost=hostname.domain.tld
StartLimitIntervalSec=0

[Service]
ExecStart=%h/YOLObotProd/Bot.sh
Restart=always
RestartSec=5
Type=simple
WorkingDirectory=%h/YOLObotProd

[Install]
WantedBy=multi-user.target
